# root/etc/services.d/flood/run
#!/usr/bin/with-contenv bash

cd /usr/flood || exit

su abc -c "flood"
