# root/etc/services.d/flood/run
#!/usr/bin/with-contenv bash

su abc -c "flood"
